# rutorrent/root/etc/services.d/rtorrent/run
#!/bin/sh
exec s6-setuidgid rtorrent /usr/local/bin/rtorrent
