# rutorrent/root/etc/services.d/rtorrent/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid rtorrent /usr/local/bin/rtorrent
