# rutorrent/root/etc/services.d/rtorrent/run
#!/usr/bin/with-contenv sh

export HOME=/home/rtorrent
exec s6-setuidgid rtorrent /usr/local/bin/rtorrent
